# deploy/supervisor.ini
[program:bilibili]
directory=/root/B/bili/bilibili-live-record-and-upload
command=/root/B/bili/bin/python3 -u run.py
stdout_logfile=/root/B/bili/bilibili-live-record-and-upload/log
redirect_stderr=true
stopwaitsecs=600
startretries=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
